# data-retrieval/app/crontab
* * * * * python /code/scripts/pop_db_test.py